# deploy/systemd/streamgui.service
[Unit]
Description=Streamdelay GUI (Flask via gunicorn)
After=network-online.target
Wants=network-online.target
# Start even if streamdelay is down; UI can still show key screen
# If you want hard dependency, uncomment:
# Requires=streamdelay.service
# After=streamdelay.service

[Service]
Type=simple
WorkingDirectory=/opt/streamgui
Environment="STREAMDELAY_API=http://127.0.0.1:7070"
Environment="STREAMDELAY_API_KEY=lan-key"
Environment="STREAMGUI_STATE_DIR=/var/lib/streamgui"
ExecStart=/opt/streamgui/.venv/bin/gunicorn -w 2 -b 127.0.0.1:8080 app:app
Restart=always
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
